# packaging/redhat/7/openpaas/openpaas.postinst
#! /bin/sh -e

setsebool -P httpd_can_network_connect 1

systemctl daemon-reload
systemctl enable openpaas
systemctl start openpaas
